# docker/root/etc/services.d/book-web/run
#!/usr/bin/with-contenv bash
exec \
    s6-setuidgid abc \
    /opt/book-web/book-web daemon -c /opt/book-web/book-web.jsonnet