# bin/book-web.service
[Unit]
Description=book-web service
After=network.target
 
[Service]
#User=king
Type=simple
ExecStart=/opt/book-web/book-web daemon
KillMode=control-group
Restart=on-failure
 
[Install]
WantedBy=multi-user.target